# backend/deploy/spending-tracker-api.service
[Unit]
Description=Spending Tracker API (NestJS)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=spending
Group=spending
WorkingDirectory=/opt/spending-tracker/backend
EnvironmentFile=/opt/spending-tracker/backend/.env
ExecStart=/usr/bin/node dist/main.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/spending-tracker/backend
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
